# Detection_Rules/Azure_Analytics/Pending_testing/3rd Party rules /BAH_Defender-TaskScheduleCreation.json
{
    "id":  "/subscriptions/41de2c2e-1ec0-4d1b-b768-edd976cdd8fe/resourceGroups/dod365-law-rg2/providers/Microsoft.OperationalInsights/workspaces/dod365-law/providers/Microsoft.SecurityInsights/alertRules/82f90c07-d334-4dea-a3ed-c3d2e4112dc9",
    "name":  "82f90c07-d334-4dea-a3ed-c3d2e4112dc9",
    "etag":  "\"87009bf2-0000-2700-0000-66fdcb9e0000\"",
    "type":  "Microsoft.SecurityInsights/alertRules",
    "kind":  "Scheduled",
    "properties":  {
                       "queryFrequency":  "PT30M",
                       "queryPeriod":  "PT30M",
                       "triggerOperator":  "GreaterThan",
                       "triggerThreshold":  0,
                       "severity":  "Medium",
                       "query":  "//Author: Steven Wan\n//Description: This analytic looks for any creation of scheduled tasks that were not completed by the system account.\n//False Positives: None known\nlet nessus_regex = \"nessus_[A-Z0-9]{8}\";\nDeviceProcessEvents\n| where FileName == \"schtasks.exe\"\n| where ActionType == \"ProcessCreated\"\n| where AccountName != \"system\"\n| where countof(InitiatingProcessCommandLine, nessus_regex, \"regex\") != 3\n| project TimeGenerated, AccountName, DeviceName, ActionType, InitiatingProcessCommandLine, InitiatingProcessFileName, InitiatingProcessParentFileName",
                       "suppressionDuration":  "PT1H",
                       "suppressionEnabled":  false,
                       "incidentConfiguration":  {
                                                     "createIncident":  true,
                                                     "groupingConfiguration":  {
                                                                                   "enabled":  false,
                                                                                   "reopenClosedIncident":  false,
                                                                                   "lookbackDuration":  "PT5M",
                                                                                   "matchingMethod":  "AllEntities",
                                                                                   "groupByEntities":  [

                                                                                                       ],
                                                                                   "groupByAlertDetails":  [

                                                                                                           ],
                                                                                   "groupByCustomDetails":  [

                                                                                                            ]
                                                                               }
                                                 },
                       "entityMappings":  [
                                              {
                                                  "entityType":  "Account",
                                                  "fieldMappings":  [
                                                                        {
                                                                            "identifier":  "AadUserId",
                                                                            "columnName":  "AccountName"
                                                                        }
                                                                    ]
                                              },
                                              {
                                                  "entityType":  "Host",
                                                  "fieldMappings":  [
                                                                        {
                                                                            "identifier":  "HostName",
                                                                            "columnName":  "DeviceName"
                                                                        }
                                                                    ]
                                              }
                                          ],
                       "eventGroupingSettings":  {
                                                     "aggregationKind":  "SingleAlert"
                                                 },
                       "tactics":  [
                                       "Persistence"
                                   ],
                       "techniques":  [
                                          "T1053"
                                      ],
                       "displayName":  "BAH_Defender-TaskScheduleCreation",
                       "enabled":  false,
                       "description":  "This analytic looks for any creation of scheduled tasks that were not completed by the system account.\nPLATFORM:Windows\nMITRE:TA0003;T1053",
                       "alertRuleTemplateName":  null,
                       "lastModifiedUtc":  "2024-10-02T22:39:25.9362803Z"
                   }
}
